# scripts/ops/quorrin-mt5-supervisor.service
[Unit]
Description=Quorrin MT5 Bridge Supervisor
Documentation=https://github.com/anomalyco/opencode
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/quorrin/mt5.env
Environment="MT5_BRIDGE_PORT=9879"
Environment="WINE_PREFIX=/root/.wine_mt5"
WorkingDirectory=/home/quorrin/quorrin
ExecStart=/home/quorrin/.venv/bin/python -m scripts.ops.mt5_bridge_supervisor \
    --bridge-host 127.0.0.1 \
    --bridge-port 9879 \
    --health-port 9880 \
    --watchdog-secs 30 \
    --max-restarts 10
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quorrin-mt5-supervisor

# Hardening — match the monitor_all script's expectations
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/quorrin/quorrin/data

[Install]
WantedBy=multi-user.target
